# deploy/systemd/ollama-docker-proxy.socket
[Unit]
Description=Private Docker bridge socket for local Ollama

[Socket]
ListenStream=172.17.0.1:11434
NoDelay=true

[Install]
WantedBy=sockets.target
